# .config/systemd/user/cliphist.service
[Unit]
Description=Clipboard history daemon (cliphist)
Documentation=man:systemd.service(5)

[Service]
Type=simple
ExecStart=/usr/bin/wl-paste --watch /usr/bin/cliphist store
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
